# setup/ipwhois.service
Description=ipwhois

[Service]
Type=simple
Restart=always
RestartSec=5s
ExecStart=/opt/bin/ip-whois

[Install]
WantedBy=multi-user.target
